# _home/.config/shell/profile.d/10_path.sh
#!/usr/bin/env sh

# =============================================================================
# To be sourced from .profile ONLY
# =============================================================================

# NOTE: PATH is being set. Be aware of what's available as needed.

export PATH="$XDG_BIN_HOME:$PATH"
