# autostart/autostart/trayer.desktop
[Desktop Entry]
Type=Application
Exec=trayer &
Name=Trayer
Comment=Show tray bar for applets/indicators
