# autostart/autostart/variety.desktop
[Desktop Entry]
Name=Variety
Type=Application
Exec=variety &
Comment=Change wallpapers automatically
